# .config/systemd/user/autotrash.service
[Unit]
Description=Empty trash
Documentation=https://github.com/bneijt/autotrash

[Service]
Type=oneshot
ExecStart="/usr/sbin/autotrash" -d 14
