# .config/systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -b -p -e -P "wall 'Charger inserted: %s - Level is %s'" -U "wall 'Charger removed: %s - Level is %s'" -f 100 -w 20 -c 10 -d 2

